# ops/ryo/systemd/symphony-engine.service.d/override.conf
[Service]
WorkingDirectory=/home/ryo/src/symphony/elixir
EnvironmentFile=
EnvironmentFile=/home/ryo/.config/symphony-ryo/symphony.env
ExecStartPre=
ExecStartPre=/home/ryo/src/symphony/ops/ryo/preflight.sh /home/ryo/src/symphony/ops/ryo/WORKFLOW.md
ExecStart=
ExecStart=/home/ryo/.local/bin/mise exec -- /home/ryo/src/symphony/elixir/bin/symphony --port 4000 --i-understand-that-this-will-be-running-without-the-usual-guardrails /home/ryo/src/symphony/ops/ryo/WORKFLOW.md
