# ops/ryo/systemd/auto-template.service.d/zz-linear-env.conf
[Service]
EnvironmentFile=/home/ryo/.config/symphony-ryo/symphony.env
ExecStart=
ExecStart=/bin/bash -lc 'export LINEAR_API_TOKEN="${LINEAR_API_TOKEN:-${LINEAR_API_KEY:-${SYMPHONY_RYO_LINEAR_TOKEN:-}}}"; exec /usr/bin/python3 /home/ryo/Github/kasotuosawari-design/auto_template/run_auto_template.py --host 127.0.0.1 --port 8787 --enable-llm'
